# server/supervisor.conf
[program:flask-blog]
command=/path/to/project/root/gunicorn.sh
autostrat=true
autorestart=true
